# ops/systemd/curriculum-ocr-reprocess-b-r3-monitor.service.d/alert-only.conf
[Unit]
# A missing run directory must fail the unit and reach OnFailure. A Condition
# would silently skip the unit and suppress the alert path.
ConditionPathIsDirectory=
OnFailure=curriculum-ocr-monitor-alert@%n.service

[Service]
SuccessExitStatus=10
EnvironmentFile=%h/.config/bdfz/curriculum-ocr-monitor-alert.conf
ExecStartPre=/usr/bin/test -d ${BDFZ_OCR_ALERT_RUN_ROOT}
ExecStartPre=/usr/bin/test -d %h/.local/state/bdfz-curriculum-ocr-monitor-alert
ExecStartPre=/usr/bin/test -f %h/curriculum-ocr-offload/alert-runtime/notify-remote-ocr-single-shard-monitor.mjs
ExecStartPost=/usr/bin/node %h/curriculum-ocr-offload/alert-runtime/notify-remote-ocr-single-shard-monitor.mjs \
  --mode observe \
  --state-dir %h/.local/state/bdfz-curriculum-ocr-monitor-alert \
  --latest-json ${BDFZ_OCR_ALERT_LATEST_JSON} \
  --expected-run-id ${BDFZ_OCR_ALERT_EXPECTED_RUN_ID} \
  --monitor-unit %n \
  --worker-unit ${BDFZ_OCR_ALERT_WORKER_UNIT} \
  --monitor-script ${BDFZ_OCR_ALERT_MONITOR_SCRIPT} \
  --monitor-sha256 ${BDFZ_OCR_ALERT_MONITOR_SHA256} \
  --max-latest-age-seconds 300
ReadWritePaths=%h/.local/state/bdfz-curriculum-ocr-monitor-alert
