# ops/launchd/com.suen.curriculum-ocr-watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.suen.curriculum-ocr-watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-dimsu</string>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/ylsuen/CF/curriculum-atlas/scripts/ocr-watchdog.mjs</string>
    <string>run</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ylsuen/CF/curriculum-atlas</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PADDLE_PDX_CACHE_HOME</key>
    <string>/Users/ylsuen/CF/curriculum-atlas/.cache/paddlex</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>ThrottleInterval</key>
  <integer>10</integer>
  <key>StandardOutPath</key>
  <string>/Users/ylsuen/CF/curriculum-atlas/.cache/ocr-supervisor/watchdog.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ylsuen/CF/curriculum-atlas/.cache/ocr-supervisor/watchdog.stderr.log</string>
</dict>
</plist>
