# ops/systemd/curriculum-ocr-reprocess-b-r3-monitor.service
[Unit]
Description=Curriculum Atlas fail-closed B-r3 parallel-1 OCR monitor
After=default.target

[Service]
Type=oneshot
UMask=0077
WorkingDirectory=/home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/workspace-b-r3
ExecStartPre=/usr/bin/test -d /home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/workspace-b-r3
ExecStartPre=/usr/bin/test -d /home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/output/production-p4-mb16-shard-b-r1
ExecStartPre=/usr/bin/test -d /home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/output/production-p1-mb16-shard-b-r3
ExecStartPre=/usr/bin/test -d /home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/monitor-b-r3
ExecStartPre=/usr/bin/sha256sum --check --strict /home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/workspace-b-r3/SHA256SUMS
ExecStart=/usr/bin/node /home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/workspace-b-r3/scripts/monitor-remote-ocr-single-shard.mjs --run-root /home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess --predecessor-output output/production-p4-mb16-shard-b-r1 --successor-output output/production-p1-mb16-shard-b-r3 --output-dir /home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/monitor-b-r3 --worker-unit curriculum-ocr-reprocess-b-r3.service --old-worker-unit a=curriculum-ocr-reprocess@a.service --old-worker-unit b=curriculum-ocr-reprocess@b.service --inactive-worker-unit b-r2=curriculum-ocr-reprocess-b-r2.service --inactive-worker-unit llama-p1-canary=curriculum-ocr-llama-p1-canary.service --llama-unit curriculum-ocr-llama.service --llama-health-url http://127.0.0.1:8112/health --b1-identity-sha256 83d9b65f772682465792ee4f76ed77d42a46541ff5e87bcfaa1559ba28a1f374 --b1-run-status-sha256 5a98d8ee9f614543ba89366e911554e31f9576c6b57953e14119563a18d6d209 --b1-state-hashset-sha256 c61c02947563caef884db15c5e60f1e79a818eaaeee0609d7502e51c0321703a --b1-status-hashset-sha256 1b41d39c07d331d35f6edd9818e02dce1bbf52c095356fd45e8471292a4f61ca --b1-artifact-hashset-sha256 3287cc5314e105dd012faadd8e22b629fd0d3f9b4fd9449aff5ccbc22b861e4c --stall-seconds 1500 --disk-min-gib 50 --memory-min-gib 2 --gpu-max-c 85
SuccessExitStatus=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadOnlyPaths=/home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess
ReadWritePaths=/home/suen/curriculum-ocr-offload/runs/20260716T1520Z-partial14-reprocess/monitor-b-r3
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
