# ops/systemd/curriculum-ocr-llama.service
[Unit]
Description=Curriculum Atlas pinned PaddleOCR-VL 1.6 llama.cpp CUDA parallel-1 server
StartLimitIntervalSec=300
StartLimitBurst=3
Conflicts=curriculum-ocr-llama-p1-canary.service
ConditionFileIsExecutable=/home/suen/curriculum-ocr-offload/runs/20260716T0250Z-paddleocrvl16-canary/src/llama.cpp/build-cuda/bin/llama-server
ConditionPathExists=/home/suen/curriculum-ocr-offload/runs/20260716T0250Z-paddleocrvl16-canary/models/PaddleOCR-VL-1.6-GGUF.gguf
ConditionPathExists=/home/suen/curriculum-ocr-offload/runs/20260716T0250Z-paddleocrvl16-canary/models/PaddleOCR-VL-1.6-GGUF-mmproj.gguf

[Service]
Type=exec
ExecStart=/home/suen/curriculum-ocr-offload/runs/20260716T0250Z-paddleocrvl16-canary/src/llama.cpp/build-cuda/bin/llama-server -m /home/suen/curriculum-ocr-offload/runs/20260716T0250Z-paddleocrvl16-canary/models/PaddleOCR-VL-1.6-GGUF.gguf --mmproj /home/suen/curriculum-ocr-offload/runs/20260716T0250Z-paddleocrvl16-canary/models/PaddleOCR-VL-1.6-GGUF-mmproj.gguf --host 127.0.0.1 --port 8112 --temp 0 --ctx-size 32768 --parallel 1 --n-gpu-layers all --mmproj-offload --flash-attn auto --cache-type-k f16 --cache-type-v f16 --batch-size 2048 --ubatch-size 512 --cont-batching --fit off --timeout 3600 --threads 8 --threads-batch 16 --no-webui --metrics
Restart=on-failure
RestartSec=5s
KillMode=control-group
TimeoutStopSec=30s
StandardOutput=append:/home/suen/curriculum-ocr-offload/runs/20260716T0250Z-paddleocrvl16-canary/logs/production-llama-p1.log
StandardError=append:/home/suen/curriculum-ocr-offload/runs/20260716T0250Z-paddleocrvl16-canary/logs/production-llama-p1.log
